# etc/systemd/system/fstrim.service.d/override.conf
[Service]
ExecStart=
ExecStart=/sbin/fstrim -v / ; /sbin/fstrim -v /boot
