# config/systemd/user/check-battery.service
[Unit]
Description=check-battery

[Service]
ExecStart=/home/luis/.config/systemd/user/check-battery.sh
